# conf/supervisor.ini
[program:celery]
command=path/bin/celery -A tasks.worker worker -Q crawl_queue -c 10 -l info -E
directory=$HOME/CrawlerMonitor/
stdout_logfile=$HOME/CrawlerMonitor/logs/celery.log
autorestart=true
redirect_stderr=true

[program:trigger]
command=path/bin/python trigger.py
directory=$HOME/CrawlerMonitor/
stdout_logfile=$HOME/CrawlerMonitor/logs/trigger.log
redirect_stderr=true
